# roles/apps/sitemcp/templates/sitemcp-cache-refresh.service.j2
[Unit]
Description=SiteMCP cache refresh with automatic suspend
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ user.name }}
Group={{ user.name }}
WorkingDirectory={{ user.home }}
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME={{ user.home }}"

# Run the cache refresh script
ExecStart={{ sitemcp_cache_refresh_script }}

# After successful completion, schedule next wake and suspend
# We use a wrapper script to handle rtcwake and suspend
ExecStartPost=/bin/bash -c 'systemctl suspend'

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sitemcp-cache-refresh

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ user.home }}/.cache/sitemcp {{ user.home }}/Documents

[Install]
WantedBy=multi-user.target
