# roles/apps/codex-cli/templates/mcp-cache-update.timer.j2
[Unit]
Description=Daily MCP cache update for Codex CLI
Documentation=man:systemd.timer(5)

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
RandomizedDelaySec=5m

# Wake system from suspend if supported by hardware
WakeSystem={{ codex_cli_cache_timer_wake_system | default(false) | bool | lower }}

[Install]
WantedBy=timers.target
