# roles/apps/codex-cli/templates/mcp-cache-update.service.j2
[Unit]
Description=Update MCP site caches for Codex CLI
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ user.name }}
Group={{ user.name }}
WorkingDirectory={{ user.home }}
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME={{ user.home }}"
Environment="MCP_CACHE_LOG_FILE=/var/log/mcp-cache/mcp-cache.log"
ExecStart={{ codex_cli_cache_script }} --all
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-cache-update

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ user.home }}/.cache/sitemcp /var/log/mcp-cache

[Install]
WantedBy=multi-user.target
